# etc/supervisor/conf.d/chat.conf
[program: chat]
command=/usr/bin/npm start
directory=/var/www/html/chat.hashbang.sh
autorestart=unexpected
stdout_logfile=/var/log/chat.stdout.log
stderr_logfile=/var/log/chat.stderr.log
user=chat
